# misc/backup/crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXCDfHrM installed on Wed Jul 20 11:29:13 2011)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
MAILTO="zsp007@gmail.com"
# m h  dom mon dow   command

#PATH42QU=/home/work/nuva/42qu
#PATH42QUMYTOOL=/home/work/nuva/42qu/mytool
PYTHON=/usr/bin/python
ZPAGE=/home/work/zpage
ZPAGE_CRONTAB=/home/work/zpage/misc/crontab
ZPAGE_BACKUP=/home/work/zpage/misc/backup

30 */3 * * * $PYTHON $ZPAGE_CRONTAB/zsite_show_buzz.py
9  2 * * * $PYTHON $ZPAGE_CRONTAB/zsite_uv.py
30 4 * * * $PYTHON $ZPAGE_CRONTAB/re_send_verify.py  
30 4 * * * $PYTHON $ZPAGE/misc/search/index.py 
30 5 * * * $PYTHON $ZPAGE/misc/sitemap/sitemap.py
3  0 * * * $PYTHON $ZPAGE_CRONTAB/notice_mail.py
9 */1 * * * $PYTHON $ZPAGE_CRONTAB/buzz_follow_mail.py
20 1 * * * $ZPAGE_BACKUP/backup_data.sh > /dev/null
1  1 * * * $ZPAGE_BACKUP/backup_crontab.sh

#30 */2 * * * $PYTHON $PATH42QUMYTOOL/rss/man_rss_update.py 
#0 */2 * * * $PYTHON $PATH42QUMYTOOL/rss/man_rss_sub.py 
#45 */3 * * * $PYTHON $PATH42QUMYTOOL/rss/man_rss_flush.py 
##4 */6 * * * $PYTHON $PATH42QUMYTOOL/rss/man_link_rss.py 
#
#30  4  * * *  $PYTHON $PATH42QU/myfile/misc/gensitemap.py
#30  6  * * * $PATH42QUMYTOOL/backup_table/backup_table.sh  > /dev/null
#
##30  5   * * * $PATH42QUMYTOOL/backup_table/backup_data.sh  > /dev/null
#30  8   * * * $PATH42QUMYTOOL/backup_table/backup_crontab.sh  > /dev/null
#
#30 4 * * *  $PYTHON $PATH42QUMYTOOL/search/index.py
#
#6 6 * * *  $PATH42QUMYTOOL/logresys/run_log2userip.sh
#
#3 3 * * *  $PYTHON $PATH42QUMYTOOL/crontab_day/mail_verify_recall.py
#1 1 * * *  $PYTHON $PATH42QUMYTOOL/crontab_day/mail_rm_man_pic.py
#7 7 * * *  $PYTHON $PATH42QUMYTOOL/crontab_day/log_history.py
#8 5 * * *  $PYTHON $PATH42QUMYTOOL/crontab_day/mail_sub_sync.py 
#30 */4 * * *  $PYTHON $PATH42QUMYTOOL/crontab_day/import_contact.py
#7 7 1 * *  $PYTHON $PATH42QUMYTOOL/crontab_day/man_ck_timeout.py
#6 4 * * *  $PYTHON $PATH42QUMYTOOL/crontab_day/man_feed_has_new.py 
#*/5 * * * *  $PYTHON $PATH42QUMYTOOL/mail_report/party_invite.py  
#3 3 * * *  $PYTHON $PATH42QUMYTOOL/mail_report/msg_expired.py  
#4 4 * * *  $PYTHON $PATH42QUMYTOOL/mail_report/msg_expire_tip.py  
#9 9 * * *  $PYTHON $PATH42QUMYTOOL/crontab_day/man_show_career_sync.py 
#3 3 * * *  $PYTHON $PATH42QUMYTOOL/crontab_day/question_sync.py 
#10 10 * * *  $PATH42QUMYTOOL/rank/man_show_rank_sort  
#
#*/15 * * * *  $PYTHON $PATH42QUMYTOOL/syncout/bind.py 
#*/15 * * * *  $PYTHON $PATH42QUMYTOOL/mail_report/namecard_apply.py
#*/15 * * * *  $PYTHON $PATH42QUMYTOOL/syncout/news.py 
#1 */6 * * *  $PYTHON $PATH42QUMYTOOL/syncout/news_follow.py
#*/20 * * * *  $PYTHON $PATH42QUMYTOOL/syncout/news_txt_db_sync.py
#
#
#*/1 * * * * $PYTHON $PATH42QUMYTOOL/crontab_min/run.py
