# linux/service/server_api_dev/run
#!/bin/bash
sudo -u work /home/work/zpage/server_api_dev.py 
